# vufind/cron.d/backups
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=$EMAIL

# Backups
5 20    * * *   root    cron-solr-backup.sh
4 20    * * *   root    cron-alpha-backup.sh
1 20    * * *   root    cron-db-backup.sh
