# vufind/cron.d/harvest
# Trigger full harvest
0 1	1 * *	root    /harvest-and-import.sh --oai-harvest --full --batch-import --yes

# Trigger incremental harvests
#TODO
