# vufind/cron.d/harvest
# Trigger full harvest
0 1	1 * *	root    /harvest.sh --full

# Trigger incremental harvests
#TODO
