# vufind/cron.d/harvest
# Trigger full harvest
0 1	1 * *	root    /harvest-and-import.sh --oai-harvest --copy-harvest --batch-import

# Trigger incremental harvests
#TODO
